# root/etc/s6-overlay/s6-rc.d/init-mod-wireguard-pia-scripts/run
../../../cont-init.d/97-create-pia-config